# nfs/etc/inittab
::sysinit:/etc/init.d/rcS
::askfirst:-/bin/sh
::ctrlaltdel:-/bin/reboot
::shutdown:/bin/umount -a -r
::restart:sbin/init
